# docker/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:crond]
command=crond -f -L 15
autostart=true

[program:queue-work]
command=php /var/www/html/artisan queue:work
autostart=true

[program:schedule-run]
command=php /var/www/html/artisan schedule:run
autostart=true
